# deploy/agent-minder-daemon.service
[Unit]
Description=agent-minder deploy daemon
Documentation=https://github.com/aptx-health/agent-minder
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=minder
Group=minder
WorkingDirectory=/opt/agent-minder/repos

# Load secrets from environment file (create from agent-minder.env.example)
EnvironmentFile=/etc/agent-minder/agent-minder.env

# Wrapper script builds CLI flags from environment variables
ExecStart=/opt/agent-minder/start-daemon.sh

Restart=on-failure
RestartSec=30
StartLimitIntervalSec=300
StartLimitBurst=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/minder/.agent-minder
ReadWritePaths=/home/minder/.claude
ReadWritePaths=/opt/agent-minder/repos
PrivateTmp=true

# Resource limits
MemoryMax=4G
MemoryHigh=3G
CPUQuota=200%
LimitNOFILE=65536
TasksMax=256

# Logging to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agent-minder

[Install]
WantedBy=multi-user.target
