# deploy/agent-minder-daemon.service
[Unit]
Description=agent-minder deploy daemon
Documentation=https://github.com/aptx-health/agent-minder
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=minder
Group=minder
WorkingDirectory=/opt/agent-minder/repo

# Load secrets from environment file (create from agent-minder.env.example)
EnvironmentFile=/etc/agent-minder/agent-minder.env

# Use the wrapper script to conditionally build --milestone or --label flag
ExecStart=/opt/agent-minder/deploy/start-daemon.sh

Restart=on-failure
RestartSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/minder/.agent-minder
ReadWritePaths=/home/minder/.claude
ReadWritePaths=/opt/agent-minder/repo
PrivateTmp=true

# Logging to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agent-minder

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
