# deploy/agent-minder-discord.service
[Unit]
Description=agent-minder Discord bot
Documentation=https://github.com/aptx-health/agent-minder
After=agent-minder-daemon.service
Wants=agent-minder-daemon.service

[Service]
Type=simple
User=minder
Group=minder

# Load secrets from environment file
EnvironmentFile=/etc/agent-minder/agent-minder.env

ExecStart=/usr/local/bin/agent-minder discord \
    --remote ${SERVE_ADDR:-localhost:7749} \
    --token "${DISCORD_BOT_TOKEN}" \
    --channel "${DISCORD_CHANNEL_ID}" \
    --guild "${DISCORD_GUILD_ID}" \
    --api-key "${MINDER_API_KEY}"

Restart=on-failure
RestartSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/minder/.agent-minder
PrivateTmp=true

# Logging to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agent-minder-discord

[Install]
WantedBy=multi-user.target
